# i486/porteus/base/003-lxqt/usr/share/gvfs/mounts/obexftp.mount
[Mount]
Type=obex
Exec=/usr/libexec/gvfsd-obexftp
AutoMount=false
