# i586/porteus/base/003-mate/usr/share/gvfs/mounts/archive.mount
[Mount]
Type=archive
Exec=/usr/libexec/gvfsd-archive
AutoMount=false
Scheme=archive
