# i586/porteus/base/003-mate/usr/share/gvfs/mounts/cdda.mount
[Mount]
Type=cdda
Exec=/usr/libexec/gvfsd-cdda
AutoMount=false
