# x86_64/porteus/base/003-xfce/usr/share/gvfs/mounts/afp-browse.mount
[Mount]
Type=afp-server
Exec=/usr/libexec/gvfsd-afp-browse
AutoMount=false
Scheme=afp
DefaultPort=548
HostnameIsInetAddress=true

