# x86_64/porteus/base/001-core/lib64/xfsprogs/xfs_scrub_all.cron
10 3 * * 0 root test -e /run/systemd/system || /sbin/xfs_scrub_all
