# i586/porteus/base/003-cinnamon/etc/xdg/autostart/cinnamon-autostart.desktop
[Desktop Entry]
Type=Application
Name=Cinnamon GUI helper
Comment=Cinnamon GUI helper
Exec=/opt/porteus-scripts/cinnamon/cinnamon-autostart

