# x86_64/porteus/base/003-kde5/etc/profile.d/sshaskpass.sh
#!/bin/bash
export SSH_ASKPASS="/usr/bin/ksshaskpass"