# i586/porteus/base/003-cinnamon/etc/profile.d/xdg.csh
#!/bin/bash

# Script to set XDG_MENU_PREFIX since it is not being set
export XDG_MENU_PREFIX=cinnamon-