# i586/porteus/base/002-xorg/usr/share/gvfs/mounts/afc.mount
[Mount]
Type=afc
Exec=/usr/libexec/gvfsd-afc
AutoMount=false
Scheme=afc
DefaultPort=1

