# i586/porteus/base/002-xorg/usr/share/gvfs/mounts/ftps.mount
[Mount]
Type=ftps
Exec=/usr/libexec/gvfsd-ftp
AutoMount=false
Scheme=ftps
DefaultPort=21
HostnameIsInetAddress=true
