# i586/porteus/base/002-xorg/usr/share/gvfs/mounts/trash.mount
[Mount]
Type=trash
Exec=/usr/libexec/gvfsd-trash
AutoMount=true

