# x86_64/porteus/base/003-kde/etc/xdg/autostart/kde-autostart.desktop
[Desktop Entry]
Type=Application
Name=KDE GUI helper
Comment=KDE GUI helper
Exec=/opt/porteus-scripts/kde/kde-autostart

