# x86_64/porteus/base/001-core/etc/cron.daily/slocate
#!/bin/sh
ionice -c3 nice -n 19 /usr/bin/updatedb -c /etc/updatedb.conf
