# i586/porteus/base/003-kde5/usr/lib/systemd/user/app-org.kde.spectacle.service
[Unit]
Description=Spectacle screenshot capture utility

[Service]
ExecStart=/usr/bin/spectacle --dbus
BusName=org.kde.Spectacle
