# i586/porteus/base/003-kde5/usr/lib/systemd/user/plasma-kded.service
[Unit]
Description=KDE Daemon
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/kded5
BusName=org.kde.kded5
Slice=session.slice
