# i586/porteus/base/003-kde5/usr/lib/systemd/user/plasma-kactivitymanagerd.service
[Unit]
Description=KActivityManager Activity manager Service
PartOf=graphical-session.target

[Service]
ExecStart=/usr/lib/kactivitymanagerd
BusName=org.kde.ActivityManager
Slice=background.slice
