# i586/porteus/base/003-mate/etc/xdg/autostart/mate-autostart.desktop
[Desktop Entry]
Type=Application
Name=MATE GUI helper
Comment=MATE GUI helper
Exec=/opt/porteus-scripts/mate/mate-autostart

