# x86_64/porteus/base/003-kde5/usr/lib/systemd/user/plasma-ksmserver.service
[Unit]
Description=KDE Session Management Server
Wants=plasma-kcminit.service
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/ksmserver
BusName=org.kde.ksmserver
Slice=session.slice

[Install]
WantedBy=plasma-core.target
