# x86_64/porteus/base/003-lxqt/usr/share/gvfs/mounts/mtp.mount
[Mount]
Type=mtp
Exec=/usr/libexec/gvfsd-mtp
AutoMount=false
