# x86_64/porteus/base/002-xorg/usr/share/gvfs/mounts/gphoto2.mount
[Mount]
Type=gphoto2
Exec=/usr/libexec/gvfsd-gphoto2
AutoMount=false
