# x86_64/porteus/base/002-xorg/usr/share/gvfs/mounts/recent.mount
[Mount]
Type=recent
Exec=/usr/libexec/gvfsd-recent
AutoMount=true

