# x86_64/porteus/base/003-mate/etc/xdg/autostart/mate-notification-daemon.desktop
[Desktop Entry]
Type=Application
Name=MATE Notification Daemon
Exec=/usr/lib64/mate-notification-daemon/mate-notification-daemon -r
OnlyShowIn=MATE;
