# x86_64/porteus/base/002-xorg/etc/xdg/autostart/gui-autostart.desktop
[Desktop Entry]
Type=Application
Name=Porteus GUI helper
Comment=Porteus GUI helper
Exec=/opt/porteus-scripts/xorg/gui-autostart

