# x86_64/porteus/base/003-openbox/etc/profile.d/ob-xdg.sh
#!/bin/bash
export XDG_DESKTOP_SESSION=openbox