# x86_64/porteus/base/003-xfce/etc/xdg/autostart/xfce4-autostart.desktop
[Desktop Entry]
Type=Application
Name=xfce4 GUI helper
Comment=xfce4 GUI helper
Exec=/opt/porteus-scripts/xfce/xfce4-autostart

